# deploy/unitybrain.service
[Unit]
Description=PinkyBrain P2P AI Network
Documentation=https://github.com/pinkybrain/pinkybrain
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=pinkybrain
Group=pinkybrain
WorkingDirectory=/opt/pinkybrain
ExecStart=/usr/local/bin/pinkybrain serve --config /etc/pinkybrain/config.json
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Environment
Environment=OLLAMA_HOST=127.0.0.1:11434
EnvironmentFile=-/etc/pinkybrain/env

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/pinkybrain /var/log/pinkybrain /run/pinkybrain
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pinkybrain

[Install]
WantedBy=multi-user.target